# server/crontab/check_inactive_users_cron
*/5 * * * * /usr/local/bin/python /app/manage.py check_inactive_users >> /var/log/cron.log 2>&1
